# config/systemd/claude-tmux.service
[Unit]
Description=Persistent Claude Code tmux session with 6-pane layout
After=network.target

[Service]
Type=forking
Environment=PATH=%h/bin:%h/.nvm/versions/node/v22.17.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/bin/claude-session -n main
ExecStop=/usr/bin/tmux kill-session -t claude-main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
